# deployment/conf/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:web]
command=/usr/local/bin/gunicorn web.wsgi:application -c /app/deployment/gunicorn.conf.py
directory=/app/
autostart=true
autorestart=true
